# etc/telex.conf
#telex bot server

description "Telex upstart script"

respawn
respawn limit 15 5

start on runlevel [2345]
stop on shutdown

setuid <username>
exec /bin/sh <telexpath>/launch.sh
